# net-vpn/nebula/files/nebula
#!/sbin/openrc-run
command="EPREFIX/usr/bin/nebula"
command_args="-config EPREFIX/etc/nebula/config.yml"
name="Nebula VPN"
description="Nebula VPN"
pidfile=EPREFIX/run/nebula.pid
command_background=true
output_log=EPREFIX/var/log/nebula.log
error_log=EPREFIX/var/log/nebula.err

depend() {
    need net
}
